# debian/syschangemon.cron.d
0 * * * * root [ -x /usr/sbin/syschangemon ] && /usr/sbin/syschangemon --skip-empty-mail 2>&1 >/dev/null
